# rc.local
nohup python3 /home/pi/Downloads/rer_capture/monitor.py &